# DEMO4/ansible/roles/role_c/files/pet.service
[Unit]
Description = set service for petclinic
[Service]
ExecStart = /home/petclinic/service.sh
[Install]
WantedBy = multi-user.target
